# apps/buddy/scripts/sara-standalone.service
[Unit]
Description=Sara Standalone Voice Loop
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/pi/reachy_mini/apps/buddy
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/pi/.buddy.env
ExecStart=/home/pi/buddy-venv/bin/python -m buddy.standalone
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
